# ObsidianAgent/Projects/IPS_IDS/adv7sec_1_0/resource_files/etc/systemd/system/loki-rs-scan.service
[Unit]
Description=IPS IDS Loki-RS YARA scan

[Service]
Type=oneshot
# [FIX-REAL-LOKI-CLI] Loki-RS 2.10.0 has no scan subcommand; -p is remote-proto.
ExecStart=/usr/local/bin/loki-rs --folder /home --no-tui
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7
