# ObsidianAgent/Projects/IPS_IDS/adv7sec_1_0/resource_files/etc/systemd/system/lynis-audit.service
[Unit]
Description=IPS IDS Lynis audit

[Service]
Type=oneshot
ExecStart=/usr/bin/lynis audit system --quiet --no-colors
