# ObsidianAgent/Projects/IC01-aops/IPS_IDS/adv7sec_1_0/resource_files/etc/audit/rules.d/50-persistence.rules
## [IPS-PH1] Workstation persistence and privilege audit rules.
## [FIX-AUDIT-APPEND] No -D/-b/-f here. This file appends IPS_IDS coverage
## without deleting or overriding audit policy owned by other hardening profiles.

## [IPS-PH1] Identity and auth surface.
-w /etc/passwd -p wa -k identity
-w /etc/shadow -p wa -k identity
-w /etc/group -p wa -k identity
-w /etc/gshadow -p wa -k identity
-w /etc/sudoers -p wa -k privilege
-w /etc/sudoers.d -p wa -k privilege
-w /etc/ssh/sshd_config -p wa -k sshd
-w /etc/ssh/sshd_config.d -p wa -k sshd
-w /home -p wa -k user-persistence

## [IPS-PH1] Persistence and system mutation.
-w /etc/systemd/system -p wa -k systemd
-w /usr/lib/systemd/system -p wa -k systemd
-w /etc/modules-load.d -p wa -k modules
-w /etc/modprobe.d -p wa -k modules
-w /etc/pacman.d/hooks -p wa -k pacman-hooks
-w /etc/nftables.conf -p wa -k firewall
-w /etc/nftables -p wa -k firewall

## [IPS-PH1] Time changes can hide activity.
-a always,exit -F arch=b64 -S adjtimex,settimeofday,clock_settime -k time-change

## [IPS-PH1] Kernel module activity.
-a always,exit -F arch=b64 -S init_module,finit_module,delete_module -k kernel-module

## [IPS-PH1] Human-session privilege escalation.
-a always,exit -F arch=b64 -S execve -F euid=0 -F auid>=1000 -F auid!=-1 -k root-cmd
-a always,exit -F arch=b32 -S execve -F euid=0 -F auid>=1000 -F auid!=-1 -k root-cmd

## [IPS-PH1] Keep mutable during tuning. Promote with -e 2 only after baselining.
