# IC01-aops/OpenBox0.1v/systemd/openbox-wg-watchdog.service
[Unit]
Description=OpenBox WireGuard handshake watchdog
After=wg-quick@wg0.service
Wants=wg-quick@wg0.service

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/openbox-wg-watchdog.sh
