# ObsidianAgent/Projects/ADV7ia/deploy/systemd/openhands-docker-proxy.service
# [FIX-OH-09] Deferred system-scope target for the OpenHands bridge proxy cutover.
[Unit]
Description=Expose OpenHands localhost API to Docker bridge only
After=network-online.target docker.service
Requires=docker.service

[Service]
Type=simple
ExecStart=/usr/bin/socat -d -d TCP4-LISTEN:3000,bind=172.17.0.1,fork,reuseaddr TCP4:127.0.0.1:3000
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
