# IC01-aops/IPS_IDS/adv7sec_1_0/resource_files/etc/systemd/system/suricata.service.d/ipsids.conf
[Service]
ExecStart=
EnvironmentFile=-/etc/default/ipsids-suricata
ExecStart=/usr/local/sbin/ipsids-suricata-run.sh
