# IC01-aops/IPS_IDS/adv7sec_1_0/resource_files/etc/systemd/system/kunai.service
[Unit]
Description=IPS IDS Kunai eBPF event monitor
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/kunai run
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target
