# IC01-aops/IPS_IDS/adv7sec_1_0/resource_files/etc/systemd/system/loki-rs-scan.timer
[Unit]
Description=Daily IPS IDS Loki-RS YARA scan

[Timer]
OnCalendar=daily
RandomizedDelaySec=45m
Persistent=true

[Install]
WantedBy=timers.target
