# ObsidianAgent/Projects/IC01-aops/IPS_IDS/adv7sec_1_0/resource_files/etc/systemd/system/aide-check.service
[Unit]
Description=IPS IDS AIDE integrity check

[Service]
Type=oneshot
ExecStart=/usr/bin/aide --check
