# ObsidianAgent/Projects/IC01-aops/OpenBox0.1v/systemd/openbox-lynis.timer
[Unit]
Description=OpenBox Lynis weekly audit (Mon 03:00)

[Timer]
OnCalendar=Mon *-*-* 03:00:00
Persistent=true
RandomizedDelaySec=30min

[Install]
WantedBy=timers.target
